# server_config/aps-docker-update.timer
[Unit]
Description=Run APS Docker Compose update nightly at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
AccuracySec=1m
Unit=aps-docker-update.service

[Install]
WantedBy=timers.target
